# ndi-to-whip.service
# =============================================================================
# /etc/systemd/system/ndi-to-whip.service
#
# Install:
#   sudo cp ndi-to-whip.service /etc/systemd/system/
#   sudo systemctl daemon-reload
#   sudo systemctl enable --now ndi-to-whip
#
# Logs:
#   journalctl -u ndi-to-whip -f
#   journalctl -u ndi-to-whip --since "1 hour ago"
# =============================================================================

[Unit]
Description=NDI to WHIP WebRTC Bridge
Documentation=https://github.com/thatdavidbelll/NDI-to-WHIP
After=network-online.target mediamtx.service
Wants=network-online.target
# Start after MediaMTX if it is installed as a service
After=mediamtx.service

# Restart limits belong here (in [Unit]), not in [Service]
StartLimitIntervalSec=300
StartLimitBurst=10

[Service]
# ── Identity ──────────────────────────────────────────────────────────────────
User=ndi-whip
Group=ndi-whip
SupplementaryGroups=video audio

# ── Command ───────────────────────────────────────────────────────────────────
ExecStart=/opt/ndi_to_whip/venv/bin/python \
    /opt/ndi_to_whip/ndi_to_whip.py \
    --config /etc/ndi_to_whip/config.toml

# ── Restart policy ────────────────────────────────────────────────────────────
Restart=on-failure
RestartSec=5s

# ── Environment ───────────────────────────────────────────────────────────────
Environment="LD_LIBRARY_PATH=/usr/local/lib"
Environment="GST_PLUGIN_PATH=/usr/local/lib/gstreamer-1.0"
Environment="NDI_SDK_DIR=/usr/local"
# Keep registry fork disabled to avoid subprocess issues under security policy
Environment="GST_REGISTRY_FORK=no"
# Point GStreamer registry cache to a writable system path
# (home dir is inaccessible under ProtectHome=true)
Environment="XDG_CACHE_HOME=/var/cache"

# ── Security hardening ────────────────────────────────────────────────────────
ReadWritePaths=/var/log/ndi_to_whip /var/cache/gstreamer-1.0
ReadOnlyPaths=/etc/ndi_to_whip
ProtectSystem=strict
ProtectHome=true
# PrivateTmp disabled: NDI SDK uses /tmp for local IPC/shared memory
NoNewPrivileges=true
# NOTE: Do NOT set CapabilityBoundingSet= (empty) — it prevents GStreamer
# codec discovery from loading plugin factories.
DeviceAllow=/dev/dri rw

# ── Resource limits ───────────────────────────────────────────────────────────
OOMScoreAdjust=-100
MemoryMax=512M

# ── Logging ───────────────────────────────────────────────────────────────────
StandardOutput=journal
StandardError=journal
SyslogIdentifier=ndi-to-whip

[Install]
WantedBy=multi-user.target
